# keepers/FetchKeepers.cron
#!/usr/bin/env bash
#
# do not modify on server! this file is under version control!
# update on development machine, commit changes, then update server
#
TIMESTAMP=`date "+%Y-%m-%d %H:%M:%S"`
LOGFILE="/tmp/searches/keeper_log.log"
if [ $(echo `date +%s` / 86400 % 2 == 1 |/usr/bin/bc) -eq 1 ]
then
        echo "$TIMESTAMP will nice FetchKeepers.sh" >> $LOGFILE
        ~/Devel/auctions/keepers/FetchKeepers.sh
else
        echo "$TIMESTAMP not today: FetchKeepers.sh" >> $LOGFILE
fi
